# crontabs/root
# SHELL=/bin/bash
# PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* * * * * root /git_pull.sh | logger
* * * * * root /git_commit.sh | logger
0 * * * * root cd /app && /bin/tar -cvjSf /tmp/"$WORLDS".tar.bz2 *orld* && /usr/bin/rclone copy /tmp/"$WORLDS".tar.bz2 "$S3_BUCKET":"$S3_BUCKET"/worlds && /bin/rm /tmp/"$WORLDS".tar.bz2 | logger
